# systemd/iptables-comment-exporter.service
[Unit]
Description=iptables-comment-exporter (/metrics по помеченным правилам iptables)
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/iptables-comment-exporter --bind=127.0.0.1:9631
Restart=on-failure
RestartSec=60

# non-root (User=prometheus + CAP_NET_ADMIN/CAP_NET_RAW через
# AmbientCapabilities) технически работает на части хостов, но ненадёжно:
# iptables-save/ip6tables-save (nftables backend) иногда отказывают в
# "fetch rule set generation id: Permission denied (you must be root)" даже
# с этими capability - без чёткой зависимости от ОС/ядра/nftables. Root -
# единственный вариант, стабильно работающий везде. См. README.md.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

# лимиты на случай зависшего/аномального *tables-save
MemoryMax=64M
CPUQuota=20%
TasksMax=16

[Install]
WantedBy=multi-user.target
